# meta-facebook/recipes-phosphor/sensors/dbus-sensors/powerOff.service
[Unit]
Description=Power off the slots when fan fails.

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/powerOff.sh
