# meta-facebook/recipes-x86/chassis/x86-power-control/yosemitev2/power_good2_asserted.service
[Unit]
Description=Power Good 2 Asserted  Service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/env toggle_power_good2.sh
