# meta-facebook/recipes-phosphor/gpio/phosphor-gpio-monitor/yosemitev2/hand_switch_position.service
[Unit]
Description=Hand Switch Position Service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/env hand_switch_position.sh  

[Install]
WantedBy=multi-user.target 
