# meta-facebook/recipes-phosphor/ipmi/phosphor-ipmi-ipmb/yosemitev2/hand_switch_position4.service
[Unit]
Description=Hand Switch Position 4 Service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/env hand_switch_position.sh  4
