# meta-facebook/recipes-phosphor/flash/phosphor-software-manager/yosemitev2/obmc-flash-host1-bios@.service
[Unit]
Description=Flash Host Bios image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/bios-update.sh host1 /tmp/images/%i
